# daemon/com.xai-personalize.dashboard.plist.tpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.xai-personalize.dashboard</string>

  <key>ProgramArguments</key>
  <array>
    <string>__INSTALL_DIR__/run.sh</string>
    <string>--no-open</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__INSTALL_DIR__</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
    <key>Crashed</key>
    <true/>
  </dict>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>/tmp/xai-personalize-dashboard.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/xai-personalize-dashboard.err</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>__HOME__/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>DASHBOARD_PORT</key>
    <string>7873</string>
  </dict>
</dict>
</plist>
